# supervisor.conf
[program:xxx]
command=xxx
process_name=%(program_name)s
directory=xxx
autostart=true
autorestart=unexpected
startsecs=10
startretries=3
user=vaayne
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
